# autossh-socks5@.service
[Unit]
Description=AutoSSH SOCKS5 Tunnel Service - %i
Documentation=man:autossh(1)
After=network-online.target
Wants=network-online.target

[Service]
# 运行环境
# AUTOSSH_GATETIME=0 禁用 AutoSSH 内置监控，依赖 SSH 层保活
Environment="AUTOSSH_GATETIME=0"
EnvironmentFile=%h/.autossh/%i.env

# 启动前清理端口（如果端口被旧的 autossh/ssh 进程占用）
ExecStartPre=/bin/bash -c 'PORT=${LOCAL_PORT}; if command -v ss &> /dev/null; then if ss -tlnp 2>/dev/null | grep -q ":$PORT "; then PID=$(ss -tlnp 2>/dev/null | grep ":$PORT " | grep -oP "pid=\\K[0-9]+" | head -1); if [ -n "$PID" ]; then PROCESS=$(ps -p $PID -o comm --no-headers 2>/dev/null || echo "unknown"); if [[ "$PROCESS" == *"autossh"* ]] || [[ "$PROCESS" == *"ssh"* ]]; then kill -15 $PID 2>/dev/null || true; sleep 2; kill -9 $PID 2>/dev/null || true; fi; fi; fi; fi; exit 0'

# 启动命令（企业级超稳定配置）
# 总超时 = 15 秒 × 120 次 = 1800 秒 = 30 分钟
# 适用于网络波动大、需要极高稳定性的场景
ExecStart=/usr/bin/autossh -M 0 -N -T \
    -i ${SSH_KEY} \
    -o ServerAliveInterval=15 \
    -o ServerAliveCountMax=120 \
    -o TCPKeepAlive=yes \
    -o ConnectTimeout=30 \
    -o ConnectionAttempts=10 \
    -D 0.0.0.0:${LOCAL_PORT} \
    ${USER}@${HOST} -p ${PORT}

# 重启策略
Restart=always
RestartSec=10

# 资源限制
LimitNOFILE=65535

# 日志输出
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autossh-%i

[Install]
WantedBy=default.target
