# autossh@.service
[Unit]
Description=AutoSSH Tunnel Service - %i
Documentation=man:autossh(1)
After=network-online.target
Wants=network-online.target

[Service]
# 运行环境
Environment="AUTOSSH_GATETIME=0"
EnvironmentFile=%h/.autossh/%i.env

# 启动前清理端口（如果端口被旧的 autossh/ssh 进程占用）
# 自动终止占用端口的旧进程，确保新服务可以启动
ExecStartPre=/bin/bash -c 'PORT=${LOCAL_PORT}; if command -v ss &> /dev/null; then if ss -tlnp 2>/dev/null | grep -q ":$PORT "; then PID=$(ss -tlnp 2>/dev/null | grep ":$PORT " | grep -oP "pid=\\K[0-9]+" | head -1); if [ -n "$PID" ]; then PROCESS=$(ps -p $PID -o comm --no-headers 2>/dev/null || echo "unknown"); if [[ "$PROCESS" == *"autossh"* ]] || [[ "$PROCESS" == *"ssh"* ]]; then kill -15 $PID 2>/dev/null || true; sleep 2; kill -9 $PID 2>/dev/null || true; fi; fi; fi; fi; exit 0'

# 启动命令
# 优化了超时配置以提高稳定性：
# - ServerAliveInterval=30: 每 30 秒发送保活消息
# - ServerAliveCountMax=10: 允许 10 次无响应（总超时 5 分钟）
# - TCPKeepAlive=yes: 启用 TCP 层保活
# - ConnectTimeout=30: 连接超时 30 秒
# - ConnectionAttempts=5: 最多尝试连接 5 次
# 注意：已移除 -C 压缩选项，减少 CPU 负担和连接不稳定
ExecStart=/usr/bin/autossh -M 0 -N -T \
    -i ${SSH_KEY} \
    -o ServerAliveInterval=30 \
    -o ServerAliveCountMax=10 \
    -o TCPKeepAlive=yes \
    -o ConnectTimeout=30 \
    -o ConnectionAttempts=5 \
    -L ${LOCAL_PORT}:localhost:${REMOTE_PORT} \
    ${USER}@${HOST} -p ${PORT}

# 重启策略
Restart=always
RestartSec=10

# 资源限制（可选）
LimitNOFILE=65535

# 日志输出
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autossh-%i

[Install]
WantedBy=default.target
